# config/systemd/user/mbsync@.service
[Unit]
Description=Mailbox %I synchronisation

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -V %i
ExecStartPost=/usr/bin/notmuch new
