# config/systemd/user/mbsync@usit.timer
[Unit]
Description=Mailbox %i timed synchronisation

[Timer]
OnBootSec=1m
OnUnitActiveSec=300

[Install]
WantedBy=timers.target

