# config/systemd/user/urlwatch.service
[Unit]
Description=Watch URLs for changes

[Service]
ExecStart=/usr/bin/urlwatch
