# config/systemd/user/dropbox.service.d/override.conf
[Service]
Environment=DISPLAY=:0
Environment="LD_PRELOAD=/usr/local/lib/libdropbox_fs_fix.so"
